# ansible/roles/uptime_kuma/templates/uptime-kuma.container.j2
[Unit]
Description=Uptime Kuma container

[Container]
Image=docker.io/louislam/uptime-kuma:{{ uptime_kuma_image_tag }}
ContainerName=uptime-kuma
PublishPort={{ uptime_kuma_publish_port }}:3001
Volume=uptime-kuma.volume:/app/data

[Service]
Restart=always
TimeoutStartSec=900

[Install]
WantedBy=multi-user.target default.target
